# docker/supervisor/env.conf
[program:lepus_proxy]
directory=/app/lepus/bin/
command=/app/lepus/bin/lepus_proxy --config=/app/lepus/etc/proxy.ini
autostart       = true
autorestart     = true
stopsignal      = TERM
redirect_stderr = true
stdout_logfile  = /dev/stdout
stdout_logfile_maxbytes = 0

[program:lepus_task]
directory=/app/lepus/bin/
command=/app/lepus/bin/lepus_task --config=/app/lepus/etc/config.ini
autostart       = true
autorestart     = true
stopsignal      = TERM
redirect_stderr = true
stdout_logfile  = /dev/stdout
stdout_logfile_maxbytes = 0

[program:lepus_alarm]
directory=/app/lepus/bin/
command=/app/lepus/bin/lepus_alarm --config=/app/lepus/etc/alarm.ini
autostart       = true
autorestart     = true
stopsignal      = TERM
redirect_stderr = true
stdout_logfile  = /dev/stdout
stdout_logfile_maxbytes = 0